# .local/share/systemd/user/auto-update.timer
[Unit]
Description=Run automatic system update weekly

[Timer]
OnCalendar=Wed
Persistent=true

[Install]
WantedBy=timers.target

