# etc/systemd/system/auto-update.service
[Unit]
Description=Automatic System Update

[Service]
Type=oneshot
ExecStart=/usr/local/bin/auto-update.sh

